# etc/profile.d/rds.sh
export PATH=$PATH:/opt/rcs-rds-linux-utils/bin

export RDS=/rds/general/user/${USER}
export RDS_PROJECT=/rds/general/project/
export EPHEMERAL=/rds/general/user/${USER}/ephemeral
